# ats/aic/resources/ansible/roles/ats/service-senza/templates/ats-senza.conf.j2
[program:ats-senza]
user = {{deploy_user}}
command = {{bin}}/senza-server
environment = {{service_environment | to_supervisor_environment}}
autostart = true
autorestart = true
stdout_logfile = {{log_dir}}/%(program_name)s.log
redirect_stderr = true
umask = 077
